# 01_Not🎃Spooky/imghost/docker/etc/supervisord.conf
[unix_http_server]
file=/dev/shm/supervisord.sock   ; (the path to the socket file)
chmod=0700                   ; socked file mode (default 0700)
chown=root

[supervisord]
logfile=/dev/shm/supervisord.log           ; (main log file;log to /dev/null, supervisor outputs to stdout)
logfile_maxbytes=0           ; (max main logfile bytes b4 rotation;0 because of /dev/null)
logfile_backups=0            ; (num of main logfile rotation backups;0 because of /dev/null)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/dev/shm/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true                ; (start in foreground if true;true because of docker)
umask=022                    ; (process file creation umask;default 022)
identifier=supervisor        ; (supervisord identifier, default is 'supervisor')

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisord.sock ; use a unix:// URL  for a unix socket

[program:php-fpm7]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/usr/sbin/php-fpm7 -F
autostart=true
autorestart=true

[program:nginx]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/usr/sbin/nginx
autostart=true
autorestart=true

